# etc/s6-overlay/s6-rc.d/dbus/run
#!/command/execlineb -P
dbus-daemon --system --nofork
