# etc/s6-overlay/s6-rc.d/add-new-groups/up
#!/command/execlineb -P
/scripts/add-new-groups < /config/groups.conf
