# etc/s6-overlay/s6-rc.d/add-new-groups/up
#!/command/execlineb -P
pipeline { cat /config/groups.conf } bash /scripts/add-new-groups.sh

