# etc/s6-overlay/s6-rc.d/wsdd/run
#!/command/execlineb -P
backtick WSDD_PARAMS { pipeline { grep WSDD_PARAMS /etc/wsdd.conf } pipeline { sed s/WSDD_PARAMS=// } sed s/\"//g }
importas wsdd_params WSDD_PARAMS
wsdd ${wsdd_params}
