# etc/s6-overlay/s6-rc.d/netatalk/run
#!/command/execlineb -P
netatalk -F /etc/netatalk/afp.conf -d
