# etc/s6-overlay/s6-rc.d/add-new-users/up
#!/command/execlineb -P
pipeline { cat /config/users.conf } bash /scripts/add-new-users.sh
