# etc/s6-overlay/s6-rc.d/add-new-users/up
#!/command/execlineb -P
/scripts/add-new-users.sh < /config/users.conf
